# service/offlineimap/run
#! /usr/bin/env bash

exec 2>&1
exec chpst -e ./env bash -c '
	exec chpst -u $USER:$GROUP offlineimap -c ./root/config
'
