# service/bitlbee/run
#! /usr/bin/env bash

exec 2>&1
exec chpst -e ./env bash -c '
	exec chpst -u $USER:$GROUP bitlbee -D -n -p 6468 -i 127.0.0.1 -d ./root
'
