# service/fossil/run
#! /usr/bin/env bash

exec 2>&1
exec chpst -e ./env bash -c '
	exec chpst -u $USER:$GROUP \
		tcpsvd -l 0 127.0.0.1 591 \
		fossil http ./root --errorlog /dev/stderr
'
